# apollo.service
[Unit]
Description=Music Recommendation App

[Service]
ExecStart=~/.local/share/apollo/apollo-server start
