# apollo.service
[Unit]
Description=Music Recommendation App

[Service]
ExecStart=/usr/sbin/apollo-server start

